# etc/service/bingebox/run
#!/usr/bin/env bash
cd /bingebox && docker-compose up